# src/urdmad/urdmad.service.in
[Unit]
Description=Create and manage virtual urdma devices

[Service]
Type=simple
ExecStart=@bindir@/urdmad --systemd
LimitMEMLOCK=infinity
RuntimeDirectory=urdma

[Install]
WantedBy=multi-user.target
